# scripts/systemd/recetas-backup-rotate.service
[Unit]
Description=Recetas - Backup rotation service
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/home/admin/api-recetas
ExecStart=/bin/bash /home/admin/api-recetas/scripts/backup_rotate.sh /home/admin/api-recetas/backups 14
User=admin
Group=admin

[Install]
WantedBy=multi-user.target
